# charms/appstream-frontend/units/sync-appstream.timer
[Unit]
Description=Sync appstream from backend (timer)

[Timer]
OnCalendar=minutely

[Install]
WantedBy=timers.target
