# charms/appstream-generator/units/appstream-generator.timer
[Unit]
Description=Appstream Generator (timer)

[Timer]
OnBootSec=5m
OnUnitInactiveSec=1h

[Install]
WantedBy=timers.target
